# scripts/core-daemon.service.in
[Unit]
Description=Common Open Research Emulator Service
After=network.target

[Service]
Type=simple
ExecStart=@PYTHON@ @bindir@/core-daemon
TasksMax=infinity

[Install]
WantedBy=multi-user.target
