# scripts/core-daemon.service.in
[Unit]
Description=Common Open Research Emulator Service
After=network.target

[Service]
Type=forking
PIDFile=/var/run/core-daemon.pid
ExecStart=@PYTHON@ @SBINDIR@/core-daemon -d
TasksMax=infinity

[Install]
WantedBy=multi-user.target
